# crontab
23 0 * * SUN /usr/local/bin/python /app/TSheetsAutomation/manage.py process_timesheets > /proc/1/fd/1 2>/proc/1/fd/2
