# crontab
23 0 * * SUN cd /home/developer/TSheets/TSheetsAutomation && /usr/local/bin/pipenv run python manage.py process_timesheets
